# crontab
# crontab
*/10 * * * * docker run --rm scraper_image